# deploy/systemd/user/kalshi-vol-surface.service
# Install: cp to ~/.config/systemd/user/ then daemon-reload. Adjust %h/kalshi-trading if needed.

[Unit]
Description=Kalshi vol-surface strategy (BTC hourly + weather)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/kalshi-trading
EnvironmentFile=-%h/kalshi-trading/.env
ExecStart=%h/kalshi-trading/.venv/bin/python -m vol_surface_strategy.main
Restart=on-failure
RestartSec=12

[Install]
WantedBy=default.target
